# user-graphical-services/compton.service
[Unit]
Description=Desktop compositor

[Service]
Type=simple
ExecStart=/usr/bin/compton

[Install]
WantedBy=user-graphical.target
