# user-graphical-services/firefox-dpms.service
[Unit]
Description=Disables DPMS when required

[Service]
Type=simple
ExecStart=/usr/bin/firefox-dpms

[Install]
WantedBy=user-graphical.target
