# user-graphical-services/labwc.service
[Unit]
Description=Labwc Wayland compositor

[Service]
Type=simple
ExecStart=/usr/bin/labwc

[Install]
WantedBy=default.target
