# user-graphical-services/mpd-custom.service
[Unit]
Description=Music Player Daemon

[Service]
Type=notify
ExecStart=/usr/bin/mpd --systemd

[Install]
WantedBy=user-graphical.target
