# user-graphical-services/mpd.service
[Unit]
Description=Music Player Daemon

[Service]
Type=notify
ExecStart=/usr/bin/mpd --no-daemon

[Install]
WantedBy=user-graphical.target
