# user-graphical-services/nitrogen.service
[Unit]
Description=Wallpaper setter

[Service]
Type=simple
ExecStart=/usr/bin/nitrogen --restore

[Install]
WantedBy=user-graphical.target
