# user-graphical-services/tint2.service
[Unit]
Description=Tint2 taskbar

[Service]
Type=simple
ExecStart=/usr/bin/tint2

[Install]
WantedBy=user-graphical.target
