# user-graphical-services/volumeicon.service
[Unit]
Description=Volume tray icon

[Service]
Type=simple
ExecStart=/usr/bin/volumeicon

[Install]
WantedBy=user-graphical.target
