# user-graphical-services/waybar.service
[Unit]
Description=Waybar taskbar

[Service]
Type=simple
ExecStart=/usr/bin/waybar

[Install]
WantedBy=user-graphical.target
