# user-graphical-services/wolumeicon.service
[Unit]
Description=Volume tray icon
Requires=panel.service
After=panel.service pipewire.service
PartOf=user-graphical.target
BindsTo=user-graphical.target

[Service]
Type=simple
ExecStart=/usr/bin/wolumeicon -m pavucontrol

[Install]
WantedBy=user-graphical.target
