# user-graphical-services/xwallpaper.service
[Unit]
Description=Wallpaper setter for X
PartOf=user-graphical.target

[Service]
Type=simple
EnvironmentFile=%h/.config/xwallpaper/config
ExecStart=/usr/bin/xwallpaper --center "$background"

[Install]
WantedBy=user-graphical.target
